# files/dropbear.conf
[program:sshd]
command=/usr/sbin/dropbear -FRs
