# files/sshd.conf
[program:sshd]
command=sh -c "ssh-keygen -A && /usr/sbin/sshd -D"
